# root/etc/s6/services/s6-syslog-ng/run
#!/usr/bin/with-contenv sh

exec syslog-ng -F --no-caps
